# tools/scripts/ralph-bridge.service
[Unit]
Description=Ralph Bridge — HTTP proxy between n8n and host AI CLIs
After=network.target
Documentation=file:///home/nntin/git/gSnake/gsnake-n8n/workflows/n8n-workflow/ralph-loop.md

[Service]
Type=simple
User=nntin
WorkingDirectory=/home/nntin/git/gSnake
EnvironmentFile=/home/nntin/git/gSnake/gsnake-n8n/tools/scripts/ralph-bridge.env
Environment=PATH=/home/nntin/.local/bin:/home/nntin/.nvm/versions/node/v20.19.6/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/nntin/.nvm/versions/node/v20.19.6/bin/node /home/nntin/git/gSnake/gsnake-n8n/tools/scripts/ralph-bridge.js
Restart=on-failure
RestartSec=5

# Logging — stdout/stderr go to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-bridge

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/nntin/git/gSnake

[Install]
WantedBy=multi-user.target
